# conf/boilerplate-supervisor.conf
[program:boilerplate]
command = /usr/bin/python /srv/run.py dev
priority = 1